# ansible/roles/orin_inference/templates/llama-server.service.j2
[Unit]
Description=llama.cpp inference server
After=network.target jetson-maxperf.service
Wants=jetson-maxperf.service

[Service]
Type=simple
User={{ llama_user }}
WorkingDirectory={{ llama_cpp_dir }}
EnvironmentFile=/etc/llama-server.env
ExecStart={{ llama_server_bin }} \
  --model {{ model_dir }}/{{ llama_default_model }} \
  --port {{ llama_port }} \
  --host 0.0.0.0 \
  --n-gpu-layers 99 \
  --ctx-size {{ llama_ctx_size }} \
  --parallel {{ llama_parallel }}
Restart=always
RestartSec=5
LimitNOFILE=65536

[Install]
WantedBy=multi-user.target
