# leds/picocrush-leds.service
[Unit]
Description=PicoCrush LED status indicator
After=network.target
Wants=llama-server.service

[Service]
Type=simple
ExecStart=/usr/bin/python3 /opt/picocluster-claw/leds/picocrush_status.py
Restart=always
RestartSec=5
User=root

[Install]
WantedBy=multi-user.target
