# leds/picocluster-claw-leds.service
[Unit]
Description=PicoCluster Claw LED status indicator
After=network.target
Wants=openclaw.service

[Service]
Type=simple
ExecStart=/usr/bin/python3 /opt/picocluster-claw/leds/picocluster_claw_status.py
Restart=always
RestartSec=5
User=root

[Install]
WantedBy=multi-user.target
